# root/lib/systemd/system/titus-setup-networking.timer
[Unit]
Description=Runs VPC Network Setup

[Timer]
OnBootSec=5s
OnUnitInactiveSec=10s

[Install]
WantedBy=timers.target
