# root/lib/systemd/system/titus-vpc-gc.timer
[Unit]
Description=Runs VPC GC

[Timer]
OnBootSec=5m
OnUnitActiveSec=3m

[Install]
WantedBy=timers.target
