# root/lib/systemd/system/titus-repear.timer
[Unit]
Description=Runs Titus Container Repear

[Timer]
OnBootSec=10m
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
