# root/lib/systemd/system/titus-vpc-gc.timer
[Unit]
Description=Runs VPC GC

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
